# src/appCataloga/root/usr/local/bin/appCataloga/run_file_bin_processing.service
[Unit]
Description=CRFS Bin file processing task service
After=network.target mariadb.service
RequiresMountsFor=/mnt/reposfi

[Service]
Type=simple
WorkingDirectory=/usr/local/bin/appCataloga
ExecStart=/usr/local/bin/appCataloga/run_file_bin_processing.sh start
ExecStop=/usr/local/bin/appCataloga/run_file_bin_processing.sh stop
ExecReload=/usr/local/bin/appCataloga/run_file_bin_processing.sh restart
Restart=on-failure

[Install]
WantedBy=multi-user.target