# src/appCataloga/root/usr/local/bin/appCataloga/run_host_task@.service
[Unit]
Description=Host task service Worker %i
After=network.target mariadb.service
RequiresMountsFor=/mnt/reposfi

[Service]
Type=simple
PIDFile=/var/run/gearman_worker_az%i.pid
WorkingDirectory=/usr/local/bin/appCataloga
ExecStart=/usr/local/bin/appCataloga/run_host_task.sh start %i
ExecStop=/usr/local/bin/appCataloga/run_host_task.sh stop %i
ExecReload=/usr/local/bin/appCataloga/run_host_task.sh restart %i
Restart=on-failure

[Install]
WantedBy=multi-user.target