# src/appCataloga/server_volume/usr/local/bin/appCataloga/OLD_Code/appCataloga_file_bkp@.service
[Unit]
Description=File backup task service Worker %i
After=network.target mariadb.service
RequiresMountsFor=/mnt/reposfi

[Service]
Type=simple
WorkingDirectory=/usr/local/bin/appCataloga
ExecStart=/usr/local/bin/appCataloga/appCataloga_file_bkp.sh start %i
ExecStop=/usr/local/bin/appCataloga/appCataloga_file_bkp.sh stop %i
ExecReload=/usr/local/bin/appCataloga/appCataloga_file_bkp.sh restart %i
Restart=on-failure

[Install]
WantedBy=multi-user.target